# roles/keitaro.roadrunner-tune/templates/systemd/roadrunner-watcher-generator.path.j2
[Path]
PathChanged={{ keitaro_watch_list }}

[Install]
WantedBy=multi-user.target
