# roles/keitaro.roadrunner-tune/templates/systemd/roadrunner-watcher.service.j2
[Unit]
Description=RoadRunner reloader
After=network.target

[Service]
Type=idle
ExecStart=/usr/local/bin/keitaro-safe_reload_roadrunner

[Install]
WantedBy=multi-user.target