# roles/keitaro.roadrunner-tune/templates/systemd/roadrunner-watcher.service.j2
[Unit]
Description=RoadRunner restarter
After=network.target

[Service]
Type=idle
ExecStart=/usr/bin/systemctl restart roadrunner.service

[Install]
WantedBy=multi-user.target