# roles/keitaro.roadrunner-tune/templates/systemd/roadrunner-watcher.path.j2
[Path]
# PathModified={{ keitaro_root }}

[Install]
WantedBy=multi-user.target
